# user/ubi2-hw21-bouygues/source-files/package/w1700k-hw21-bouygues-support/files/etc/uci-defaults/99-w1700k-hw21-bouygues-v2
#!/bin/sh

set -u

append_sysupgrade_path() {
	path="$1"
	grep -qxF "$path" /etc/sysupgrade.conf 2>/dev/null ||
		printf '%s\n' "$path" >> /etc/sysupgrade.conf
}

restore_managed_file() {
	path="$1"
	mode="$2"
	rom_path="/rom$path"

	[ -f "$rom_path" ] || return 0
	mkdir -p "${path%/*}"
	cp -f "$rom_path" "$path"
	chmod "$mode" "$path"
}

# These files are kept across generic sysupgrades so they can reinstall
# dnsproxy/full wpad. A preserved older copy otherwise masks an updated copy
# shipped by this profile, so make the current image authoritative here.
restore_managed_file /usr/local/sbin/restore-dnsproxy-after-sysupgrade 0755
restore_managed_file /etc/init.d/dnsproxy-restore 0755
restore_managed_file /etc/hotplug.d/iface/95-dnsproxy-restore 0755

append_sysupgrade_path /usr/local/sbin/restore-dnsproxy-after-sysupgrade
append_sysupgrade_path /etc/init.d/dnsproxy-restore
append_sysupgrade_path /etc/hotplug.d/iface/95-dnsproxy-restore
append_sysupgrade_path /etc/rc.d/S99dnsproxy-restore

# Root-cause patch 971 restores the stock dual-stack fw4 rule. Stop
# carrying a package-owned template into later, potentially incompatible fw4
# versions once that stock rule is present.
fw4_template=/usr/share/firewall4/templates/ruleset.uc
rom_fw4_template="/rom$fw4_template"
if [ -f /etc/sysupgrade.conf ] &&
	grep -qxF "$fw4_template" /etc/sysupgrade.conf 2>/dev/null &&
	grep -Fq 'meta l4proto { tcp, udp } flow offload @ft;' "$fw4_template" 2>/dev/null &&
	! grep -Fq 'meta nfproto ipv4 meta l4proto' "$fw4_template" 2>/dev/null &&
	grep -Fq 'meta l4proto { tcp, udp } flow offload @ft;' "$rom_fw4_template" 2>/dev/null &&
	! grep -Fq 'meta nfproto ipv4 meta l4proto' "$rom_fw4_template" 2>/dev/null; then
	# The current upgrade may already have restored an older package-owned
	# template. Refresh it now; the sysupgrade.conf edit applies only later.
	cp -f "$rom_fw4_template" "$fw4_template"
	chmod 0644 "$fw4_template"
	tmp_sysupgrade_conf="/tmp/sysupgrade.conf.$$"
	grep -vxF "$fw4_template" /etc/sysupgrade.conf > "$tmp_sysupgrade_conf" || true
	cat "$tmp_sysupgrade_conf" > /etc/sysupgrade.conf
	rm -f "$tmp_sysupgrade_conf"
fi

/etc/init.d/dnsproxy-restore disable
/etc/init.d/dnsproxy-restore enable

configure_irqbalance() {
	section=irqbalance
	[ "$(uci -q get irqbalance.$section)" = irqbalance ] ||
		uci set irqbalance.$section=irqbalance

	# Migrate only the untouched package default. A custom IRQ policy is left
	# exactly as the administrator configured it.
	if [ -z "$(uci -q get irqbalance.$section.w1700k_v2_managed)" ] &&
		[ "$(uci -q get irqbalance.$section.enabled)" = 0 ] &&
		[ -z "$(uci -q get irqbalance.$section.deepestcache)" ] &&
		[ -z "$(uci -q get irqbalance.$section.interval)" ] &&
		[ -z "$(uci -q get irqbalance.$section.banned_cpulist)" ] &&
		[ -z "$(uci -q get irqbalance.$section.banirq)" ]; then
		uci set irqbalance.$section.enabled='1'
		uci set irqbalance.$section.deepestcache='2'
		uci set irqbalance.$section.interval='10'
		uci set irqbalance.$section.w1700k_v2_managed='1'
		uci commit irqbalance
	fi

	/etc/init.d/irqbalance enable
}

watchcat_is_stock() {
	sections="$(uci -q show watchcat 2>/dev/null |
		sed -n 's/^watchcat\.\([^.=]*\)=watchcat$/\1/p')"
	[ "$(printf '%s\n' "$sections" | awk 'NF { count++ } END { print count + 0 }')" -eq 1 ] ||
		return 1
	section="$sections"
	[ "$(uci -q get "watchcat.$section.mode")" = ping_reboot ] &&
		[ "$(uci -q get "watchcat.$section.period")" = 6h ] &&
		[ "$(uci -q get "watchcat.$section.pinghosts")" = 8.8.8.8 ] &&
		[ "$(uci -q get "watchcat.$section.forcedelay")" = 30 ] &&
		[ -z "$(uci -q get "watchcat.$section.interface")" ] &&
		[ -z "$(uci -q get "watchcat.$section.script")" ] &&
		[ -z "$(uci -q get "watchcat.$section.addressfamily")" ]
}

configure_watchcat() {
	if watchcat_is_stock; then
		uci -q delete "watchcat.$section"
		uci -q batch <<'EOF'
set watchcat.wan_ipv4=watchcat
set watchcat.wan_ipv4.mode='run_script'
set watchcat.wan_ipv4.period='15m'
set watchcat.wan_ipv4.pingperiod='60s'
set watchcat.wan_ipv4.pingsize='standard'
set watchcat.wan_ipv4.pinghosts='1.1.1.1 9.9.9.9 185.222.222.222'
set watchcat.wan_ipv4.interface='@wan'
set watchcat.wan_ipv4.addressfamily='ipv4'
set watchcat.wan_ipv4.script='/usr/local/sbin/w1700k-watchcat-recover'
set watchcat.wan_ipv4.reset_failure_timer='1'
set watchcat.wan_ipv4.w1700k_v2_managed='1'
set watchcat.wan_ipv6=watchcat
set watchcat.wan_ipv6.mode='run_script'
set watchcat.wan_ipv6.period='15m'
set watchcat.wan_ipv6.pingperiod='60s'
set watchcat.wan_ipv6.pingsize='standard'
set watchcat.wan_ipv6.pinghosts='2606:4700:4700::1111 2620:fe::fe 2a09::'
set watchcat.wan_ipv6.interface='@wan6'
set watchcat.wan_ipv6.addressfamily='ipv6'
set watchcat.wan_ipv6.script='/usr/local/sbin/w1700k-watchcat-recover'
set watchcat.wan_ipv6.reset_failure_timer='1'
set watchcat.wan_ipv6.w1700k_v2_managed='1'
EOF
		uci commit watchcat
	fi

	/etc/init.d/watchcat enable
}

configure_profile_defaults() {
	if [ "$(uci -q get system.@system[0].w1700k_v2_initialized)" != 1 ]; then
		uci set firewall.@defaults[0].flow_offloading='1'
		uci set firewall.@defaults[0].flow_offloading_hw='1'
		uci commit firewall

		# Keep packet steering off for the first controlled irqbalance test.
		uci set network.globals.packet_steering='0'
		uci commit network

		uci set system.@system[0].w1700k_v2_initialized='1'
		uci commit system
	fi
}

configure_profile_defaults
configure_irqbalance
configure_watchcat

[ -x /etc/init.d/fan ] && /etc/init.d/fan enable
[ -x /etc/init.d/w1700k-pstore ] && /etc/init.d/w1700k-pstore enable

exit 0
