# user/ubi2-hw21-bouygues/source-files/package/w1700k-hw21-bouygues-support/files/etc/hotplug.d/net/50-flow-offload-wifi
#!/bin/sh

# Reload firewall after WiFi interfaces come up so fw4 can discover
# them and add them to the flowtable for hardware bridge offload via PPE.
#
# fw4 auto-discovers bridge member devices for the flowtable, but
# at boot the firewall starts before hostapd creates WiFi interfaces.

[ "$ACTION" = "add" ] || exit 0

DEV="${DEVICENAME:-$INTERFACE}"
[ -n "$DEV" ] || exit 0

case "$DEV" in
	phy*-ap*) ;;
	*) exit 0 ;;
esac

# Skip if the device is already in the flowtable (fw4 picked it up).
nft list flowtable inet fw4 ft 2>/dev/null | grep -qw "\"$DEV\"" && exit 0

# fw4 auto-discovers bridge members on reload. Serialize via flock,
# WAITING for any in-flight reload rather than skipping: a reload that
# started before this device existed snapshots the device list without
# it, so bailing out here leaves the device unoffloaded until the next.
(
	mkdir -p /var/lock
	exec 9>/var/lock/flow-offload-wifi.lock
	flock 9
	# Re-check after acquiring: the reload we waited on may have
	# already picked this device up.
	nft list flowtable inet fw4 ft 2>/dev/null | grep -qw "\"$DEV\"" && exit 0
	sleep 2
	/etc/init.d/firewall reload >/dev/null 2>&1
	if nft list flowtable inet fw4 ft 2>/dev/null | grep -qw "\"$DEV\""; then
		logger -t flow-offload-wifi "fw4 reloaded, $DEV added to flowtable"
	else
		logger -t flow-offload-wifi "WARNING: $DEV still missing from flowtable after fw4 reload"
	fi
) &

exit 0
