# user/ubi2-hw21-bouygues/source-files/package/w1700k-hw21-bouygues-support/files/etc/init.d/w1700k-pstore
#!/bin/sh /etc/rc.common

# rc.common consumes these globals after sourcing this file.
# shellcheck disable=SC2034
START=12
# shellcheck disable=SC2034
STOP=90

start() {
	/usr/local/sbin/w1700k-pstore-collect --boot >/dev/null 2>&1 || true
}
